# supervisord.conf
[supervisord]
nodaemon=true

[program:MySQL2Redis]
command=sbt run
directory=/srv/MySQL2Redis

[program:redis]
command=redis-server
directory=/srv
